# deployment/sheets-sync.service
[Unit]
Description=Google Sheets to MongoDB Sync Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/GuestListScripts
ExecStart=/usr/bin/python3 /home/ec2-user/GuestListScripts/sheets_to_mongo_sync.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/ec2-user/GuestListScripts/logs/sheets_sync_service.log
StandardError=append:/home/ec2-user/GuestListScripts/logs/sheets_sync_service.log
Environment=PYTHONPATH=/home/ec2-user/GuestListScripts

[Install]
WantedBy=multi-user.target
